# gw2/etc/local.d/99-batman-hop-penalty.start
echo 60 > /sys/class/net/ggrzBAT/mesh/hop_penalty
